# all-on-one/resources/services.d/ferretdb/run
#!/usr/bin/with-contenv sh

cd /nightscout-ferretdb
./ferretdb